# provisioner/boottime.service
[Service]
ExecStart=/root/boottime.sh
[Install]
WantedBy=default.target